# templates/com.sessionlauncher.watchdog.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd agent: runs the ttyd file-descriptor watchdog every 5 minutes.
     fd-watchdog.sh reaps leaked `tmux attach` orphans and recycles ttyd before
     it exhausts its descriptors (see the script header for the why). This is a
     periodic job, not a daemon, so it uses StartInterval and has no KeepAlive.
     The literal __HOME__ placeholder is sed-replaced with the installing user's
     home directory, so this template stays username-generic. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.sessionlauncher.watchdog</string>
  <key>ProgramArguments</key>
  <array>
    <string>__HOME__/.claude-launcher/bin/fd-watchdog.sh</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>StartInterval</key><integer>300</integer>
  <key>StandardOutPath</key><string>__HOME__/.claude-launcher/watchdog.log</string>
  <key>StandardErrorPath</key><string>__HOME__/.claude-launcher/watchdog.log</string>
</dict>
</plist>
