# templates/session-watchdog.timer
# systemd user timer: fires session-watchdog.service every 5 minutes to keep
# ttyd from leaking itself into its file-descriptor ceiling (see the service and
# fd-watchdog.sh for the why). Install under ~/.config/systemd/user/ and enable
# with `systemctl --user enable --now session-watchdog.timer`.
[Unit]
Description=Run the Session Launcher ttyd fd watchdog every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=session-watchdog.service

[Install]
WantedBy=timers.target
