# templates/session-watchdog.service
# systemd user unit: runs the ttyd file-descriptor watchdog once.
# Paired with session-watchdog.timer, which fires it every 5 minutes. Install
# under ~/.config/systemd/user/ and enable the .timer (not this unit) with
# `systemctl --user`. %h expands to the user's home directory, so this unit
# stays username-generic.
[Unit]
Description=Session Launcher ttyd fd watchdog

[Service]
Type=oneshot
ExecStart=%h/.claude-launcher/bin/fd-watchdog.sh
